# configs/systemd/news-harness-site.service
[Unit]
Description=News Harness website and read-only API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/news_harness
EnvironmentFile=/run/news-harness/news_harness.env
ExecStart=/usr/bin/python3 -m news_harness serve \
    --host 127.0.0.1 \
    --port 8765 \
    --feed web/data/radar-timeline/timeline_feed.json \
    --artifact-dir artifacts/manual_smoke/latest
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
