# configs/systemd/news-harness-cycle.timer
[Unit]
Description=Run News Harness rolling cycle every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=30s
Persistent=true
Unit=news-harness-cycle.service

[Install]
WantedBy=timers.target
